# data/deb-scripts/prerm
#!/bin/sh

systemctl stop dst-master.service || true
systemctl stop dst-caves.service || true

systemctl disable dst-master.service || true
systemctl disable dst-caves.service || true
