# deb-scripts/postinst
#!/bin/sh

echo "MyDediServer" > /etc/dst/cluster

systemctl enable dst-master.service
systemctl enable dst-caves.service

systemctl start dst-master.service
systemctl start dst-caves.service
